# systemd/ccs-plugin-cd.service
[Unit]
Description=Refresh minh-internal Claude Code plugins from their canonical repos

[Service]
Type=oneshot
Environment=PATH=/home/minh/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/mnt/2TB/minh-archive/claude-code-settings/bin/plugin-autoupdate.sh
