# systemd/ccs-plugin-cd.timer
[Unit]
Description=Schedule the minh-internal plugin CD watcher (poll canonical repos)

[Timer]
OnCalendar=*:3/10
Persistent=true

[Install]
WantedBy=timers.target
